# files/php-fpm.ini
[program:php-fpm]
command=/usr/bin/pidproxy /var/run/php-fpm/php-fpm.pid /usr/sbin/php-fpm --nodaemonize -R --allow-to-run-as-root
autostart = true
